# service_files/phal-test-app@.service.in
[Unit]
Description=phal test application

[Service]
RemainAfterExit=YES
Type=oneshot
ExecStart=/usr/bin/phal-test-app
