# service_files/phal-test-app@.service.in
[Unit]
Description=phal test application

[Service]
SyslogIdentifier=phal-test-app
RemainAfterExit=no
Type=oneshot
ExecStart=/usr/bin/phal-test-app

[Install]
WantedBy=multi-user.target
